# deploy/retrieval_node-backup.timer
[Unit]
Description=Nightly retrieval_node Postgres backup timer

[Timer]
OnCalendar=*-*-* 03:15:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
